# doc/systemd-services/update-doc-icelk.service
[Unit]
Description=Update Icelk docs

[Service]
Type=oneshot
ExecStart=bash -c "cd ~/kvarn/icelk.dev/doc && ./rebuild-docs.fish"
